# ops/systemd/ao-metrics-seed.service
[Unit]
Description=Ensure AO metrics prom file exists
After=network.target

[Service]
Type=oneshot
UMask=027
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ReadWritePaths=/var/lib/ao
ExecStart=/usr/bin/env bash -c "install -o ao -g ao -m 640 /dev/null /var/lib/ao/metrics.prom"

[Install]
WantedBy=multi-user.target
