# ops/systemd/ao-metrics-seed.service
[Unit]
Description=Ensure AO metrics prom file exists
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/env bash -c "install -o ao -g ao -m 664 /dev/null /var/lib/ao/metrics.prom"

[Install]
WantedBy=multi-user.target
